# src/modules/openclaw/filesystem/root/etc/systemd/system/openclaw.service
[Unit]
Description=OpenClaw AI Agent
After=network-online.target clawbot-swap.service
Wants=network-online.target
# Only one agent should be active at a time
Conflicts=picoclaw.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/openclaw serve
Restart=on-failure
RestartSec=10
# Higher memory allowance for Node.js
MemoryMax=600M
MemoryHigh=400M
Environment=HOME=/home/pi
Environment=NODE_OPTIONS=--max-old-space-size=384

[Install]
WantedBy=multi-user.target
