# src/modules/clawbot/filesystem/root/etc/systemd/system/clawbot-status-api.service
[Unit]
Description=ClawbotOS System Status API
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/clawbot-status-api
Restart=always
RestartSec=5
# Run as root to allow systemctl actions and reading /proc
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
