# src/modules/clawbot_wizard/filesystem/root/etc/systemd/system/clawbot-firstboot.service
[Unit]
Description=ClawbotOS First Boot Configuration
After=network-online.target clawbot-cloud.service
ConditionPathExists=/.firstboot_clawbot

[Service]
# cage manages its own VT allocation via LIBSEAT_BACKEND=direct
# DO NOT set StandardInput=tty here — it causes SIGHUP to kill the process
# when cage takes over the display (TTY session leader conflict)
ExecStart=/usr/local/bin/clawbot-firstboot
Restart=no

[Install]
WantedBy=multi-user.target
