# src/modules/smartpad/filesystem/root/etc/systemd/system/clawbot-kiosk.service
[Unit]
Description=ClawbotOS Kiosk Browser (SmartPad Dashboard)
Documentation=https://github.com/Yumi-Lab/ClawBot-OS
After=network-online.target nginx.service picoclaw.service clawbot-core.service seatd.service
Wants=network-online.target nginx.service seatd.service
# Only run on devices with a screen (SmartPad)
ConditionPathExists=/dev/dri/card0
# Only run after firstboot wizard has completed (sentinel file removed)
ConditionPathNotExists=/.firstboot_clawbot

[Service]
Type=simple
# cage refuses to run as root — run as pi (video group → seatd access)
User=pi
Group=pi
# LIBSEAT_BACKEND=seatd: seatd daemon handles DRM seat (direct not compiled in Debian)
Environment=XDG_RUNTIME_DIR=/run/clawbot-kiosk
Environment=LIBSEAT_BACKEND=seatd
Environment=HOME=/home/pi
RuntimeDirectory=clawbot-kiosk
RuntimeDirectoryMode=0700
ExecStart=/usr/bin/cage -- /usr/local/bin/clawbot-kiosk-launcher
Restart=on-failure
RestartSec=5
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
