# src/modules/armbian_net/filesystem/root/etc/systemd/system/wifi-watchdog.timer
[Unit]
Description=WiFi Watchdog Timer

[Timer]
OnBootSec=90
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
