# src/modules/armbian_net/filesystem/root/etc/systemd/system/wifi-watchdog.service
[Unit]
Description=WiFi Watchdog — reconnect on connectivity loss
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi-watchdog
