# src/modules/picoclaw/filesystem/root/etc/systemd/system/picoclaw.service
[Unit]
Description=PicoClaw AI Gateway
After=network.target clawbot-swap.service
ConditionPathExists=/home/pi/.picoclaw/config.json

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi
ExecStart=/usr/local/bin/picoclaw gateway
Restart=on-failure
RestartSec=5
# Memory limit to protect 1GB system
MemoryMax=50M
MemoryHigh=30M
Environment=HOME=/home/pi

[Install]
WantedBy=multi-user.target
