# src/modules/clawbot_core/filesystem/root/etc/systemd/system/clawbot-heartbeat.service
[Unit]
Description=ClawbotOS Cloud Heartbeat
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clawbot-heartbeat
StandardOutput=journal
StandardError=journal

# No restart — the timer drives recurrence
