# src/modules/clawbot_core/filesystem/root/etc/systemd/system/clawbot-cloud.service
[Unit]
Description=ClawbotOS Cloud WebSocket Connection
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/clawbot-cloud
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
